# net-vpn/netmaker-netclient-bin/files/netclient.initd
#!/sbin/openrc-run
# Copyright 2022 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="netclient daemon"
description=""
command=/usr/bin/netclient
command_args="daemon"
command_background=true
